# cron.d/picks
*/15 * * * * curl localhost/index.php/update
