# cron.d/picks
*/15 * * * * root curl localhost/index.php/update
